# systemd/blueman-applet.service
# ~/.config/systemd/user/blueman-applet.service
[Unit]
Description=Blueman Bluetooth Applet
After=hyprland-session.target
PartOf=hyprland-session.target

[Service]
ExecStart=/usr/bin/blueman-applet
Restart=on-failure

[Install]
WantedBy=hyprland-session.target

