# systemd/user/hypridle.service
[Unit]
Description=Hypridle - Idle daemon for Hyprland
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/hypridle
Restart=on-failure
RestartSec=2
KillMode=process
Environment=HYPRIDLE_CONFIG=%h/.config/hypr/hypridle.conf

[Install]
WantedBy=graphical-session.target

